# build/supervisor/laravel-echo-server.conf
[supervisord]
nodaemon=true

[program:laravel-echo-server]
directory=/var/www
command=/usr/bin/laravel-echo-server start
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-echo-server.log
