# build/supervisor/nginx.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=900
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log
username=www-data
autorestart=true
