# build/laravel-echo-server.conf
[program:laravel-echo-server]
directory=/var/www
command=/usr/bin/laravel-echo-server start
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-echo-server.log
